# tmp/droid-config.files
%defattr(-,root,root,-)
/lib/udev/rules.d/997-droid-bsg.rules
/lib/udev/rules.d/998-droid-system.rules
/lib/udev/rules.d/997-droid-drm.rules
/etc/xdg/QtProject/QtQuick2.conf
/etc/xdg/QtProject/Sensors.conf
/etc/xdg/QtProject/QPlatformTheme.conf
/etc/media_codecs_ffmpeg.xml
/etc/audio_platform_info.xml
/etc/media_codecs.xml
/etc/sec_config
/etc/media_codecs_performance.xml
/etc/sap.conf
/etc/gps.conf
/etc/acdbdata
/etc/media_codecs_google_video.xml
/etc/mixer_paths.xml
/etc/xtra_root_cert.pem
/etc/media_profiles.xml
/etc/obexd/opp.conf
/etc/obexd/noplugins/mas
/etc/obexd/noplugins/irmc
/etc/obexd/noplugins/syncevolution
/etc/obexd/noplugins/ftp
/etc/obexd/excludes/bluetooth:pcsuite
/etc/obexd/pbap.conf
/etc/obexd/root
/etc/lowi.conf
/etc/media_codecs_google_telephony.xml
/etc/media_codecs_google_audio.xml
/etc/selinux/config
/etc/sensorfw/primaryuse.conf
/etc/ofono/binder.d/dual-sim.conf
/etc/ofono/binder.conf
/etc/izat.conf
/etc/usb-moded/dyn-modes/mtp_mode-droid.ini
/etc/modprobe.d/blacklist-evbug.conf
/etc/xtwifi.conf
/var/lib/environment/nemo/60-multimedia.conf
/var/lib/environment/ofono/noplugin.conf
/var/lib/environment/compositor/droid-hal-device.conf
/usr/bin/droid/droid-load-firmware.sh
/usr/bin/droid/droid-hal-startup.sh
/usr/bin/droid/get_partition.sh
/usr/bin/droid/droid-hal-shutdown.sh
/usr/bin/droid/droid-reclaim-memory.sh
/usr/bin/droid/droid-init-done.sh
/usr/lib/tmpfiles.d/droid-hal.conf
/usr/lib/systemd/user/lipstick.service.d/50-compositor.conf
/usr/lib/systemd/user/jolla-startupwizard-pre-user-session.service.d/50-compositor.conf
/usr/lib/systemd/user/jolla-actdead-charging.service.d/50-compositor.conf
/usr/lib/systemd/system/droid-reclaim-memory.timer
/usr/lib/systemd/system/graphical.target.wants/droid-late-start.service
/usr/lib/systemd/system/graphical.target.wants/droid-hal-init.service
/usr/lib/systemd/system/droid-late-start.service
/usr/lib/systemd/system/timers.target.wants/droid-reclaim-memory.timer
/usr/lib/systemd/system/droid-reclaim-memory.service
/usr/lib/systemd/system/actdead.target.wants/droid-hal-init.service
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/sysctl.d/ipv4-settings.conf
/usr/lib/sysctl.d/ipv6-settings.conf
/boot/flashing-README.txt
/boot/extracting-README.txt
/etc/gbinder.conf
/etc/gst-droid/gstdroidcodec.conf
/etc/sailfish-fpd/50-settings.ini
/etc/mce/60-hybris-features.ini
/etc/mce/60-proximity-sensor.conf
/etc/mce/60-fingerprint-sensor.conf
/etc/ofono/ril_subscription.conf
/etc/ofono/ril_subscription.d/dual-sim.conf
/var/lib/environment/nemo/99-qtscenegraph.conf
/usr/share/ngfd/plugins.d/60-profile.ini
/usr/bin/droid/droid-bootctl.sh
/usr/lib/systemd/system/wait_for_keymaster.service
/usr/lib/systemd/system/graphical.target.wants/droid-bootctl.service
/usr/lib/systemd/system/droid-bootctl.service
/usr/lib/systemd/system/multi-user.target.wants/wait_for_keymaster.service
/usr/lib/systemd/system/droid-hal-init.service
/usr/libexec/droid-hybris/system/etc/init/ecclist.rc
/usr/libexec/droid-hybris/system/etc/init/disabled_services.rc
/etc/libncicore.conf
/etc/gbinder.d/anbox-hybris.conf
/etc/audio_policy_configuration.xml
/etc/mce/70-backlight.ini
/etc/selinux/minimum/contexts/files/file_contexts
/etc/selinux/minimum/contexts/dbus_contexts
/etc/selinux/minimum/policy/policy.30
/etc/selinux/config
/etc/sensorfw/sensord.conf.d/60-sensor-selection-overrides.conf
/etc/usb-moded/99-usb-moded-configfs-cepheus.ini
/etc/usb-moded/80-use-ask-mode.ini
/etc/usb-moded/run/down-usb0.ini.txt
/odm/lib64/vendor.lineage.camera.motor@1.0.so
/odm/lib64/libnativehelper.so
/odm/lib64/libandroidicu.so
/odm/lib64/libicuuc.so
/odm/lib64/libicui18n.so
/odm/lib/vendor.lineage.camera.motor@1.0.so
/odm/lib/libnativehelper.so
/odm/lib/libandroidicu.so
/odm/lib/libicuuc.so
/odm/lib/libicui18n.so
/usr/share/ngfd/plugins.d/51-ffmemless.ini
/usr/share/csd/settings.d/hw-settings.ini
/usr/bin/droid/droid-early-init.sh
/usr/lib/systemd/system/multi-user.target.wants/wlan-module-load.service
/usr/lib/systemd/system/wlan-module-load.service
/usr/libexec/droid-hybris/system/etc/init/android.hardware.media.omx@1.0-service.rc
/usr/libexec/droid-hybris/system/etc/init/init_device.rc
/usr/libexec/droid-hybris/system/etc/init/overrides.rc
/usr/libexec/droid-hybris/system/etc/init/disabled_services_device.rc
/usr/share/ssu/features.d/adaptation-community.ini
